# rootfs/etc/s6-overlay/s6-rc.d/dbus/run
#!/bin/bash

exec s6-notifyoncheck -c /etc/s6-overlay/s6-rc.d/dbus/check_readiness.sh dbus-daemon --system --nofork --nosyslog
