# rootfs/etc/s6-overlay/s6-rc.d/dbus/run
#!/bin/bash

exec dbus-daemon --system --nofork --nosyslog
